# .config/systemd/user/url-listener.service
[Unit]
Description=URL listener

[Service]
ExecStart=%h/.config/newsboat/listener.sh
Restart=always
